# src/etc/ssh/remote_install/rc.local
if [ -e /etc/ssh/remote_install/remote_install.conf ]; then
	echo -n 'sshd remote install'; /usr/sbin/sshd -f /etc/ssh/remote_install/remote_install.conf
fi

if [ -e /var/privoxy/etc/privoxy/config ]; then
        echo -n 'privoxy'; /usr/local/sbin/privoxy --user _privoxy --chroot --no-daemon --pre-chroot-nslookup www.google.com /etc/privoxy/config &
fi
